# DinD-debian-latest/supervisor/dockerd.conf
[program:dockerd]
command=/usr/local/bin/dockerd
autorestart=true
autostart=true
stdout_logfile=/var/log/dockerd.out.log
stderr_logfile=/var/log/dockerd.err.log